# com.energynewsbot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.energynewsbot</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/kim-minyeob/.nvm/versions/node/v24.14.0/bin/node</string>
        <string>/Users/kim-minyeob/.nvm/versions/node/v24.14.0/bin/npx</string>
        <string>ts-node</string>
        <string>/Users/kim-minyeob/energy-news-bot/src/index.ts</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/kim-minyeob/energy-news-bot</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/kim-minyeob/energy-news-bot/bot.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/kim-minyeob/energy-news-bot/bot-error.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/kim-minyeob/.nvm/versions/node/v24.14.0/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
